# crontab
* * * * * /app/cron-task.sh
# Mandatory blank line